# 1.1/usb-log-saver.rules
ACTION=="add", KERNEL=="sd[a-z]*", ENV{ID_BUS}=="usb", RUN+="/usr/bin/usb-added.sh"
